# zk_arcade.service
[Unit]
Description=ZkArcade
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/app/zk_arcade/
EnvironmentFile=/home/app/config/.env.zk_arcade
ExecStart=/home/app/zk_arcade/web/_build/prod/rel/zk_arcade/bin/zk_arcade start
Restart=always
RestartSec=1
StartLimitBurst=100
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target

"0x6729b62aF0B44e03f16d30C0Cb467295B2CE5A63",
"0x04301F5768A126F726cdB9e873E03549e723c988",
"0xd9327b1b9934c09a650a8ebd36dA0b17ac7D9743",
"0xD9309b3CF1B7DBF59f53461c2a66e2783dD1766f",
"0xEA4e13a1107b002aa7547dF6681b80ECC191B13E",
"0x89629A4535aE7F1693755e2e444eCF837391aaED",
"0x547045B40B3c78AD786e743F72983c304ab3aaE8",
"0xEaC5119a52eF90684d5153D2854c098cd34D5099",
"0x751Aa4C8404338B77b90Fab96e5BCdA6642270f6",
"0xB5b8F14c258F1e2648848Ed799123735D11164Ae",
"0xa06E0a632Fa4BAD7A7619e0fddeBBc39f899dA78",
"0x6F153CA3b88846B2235d05c259DbAfF28029d6A1",
"0x488976ae982F9f6928D0FA2c271253D61eA6A82b",
"0xF2533E95195E7f54220e6d56e06645eACA5d7e47",
"0xe93271FDD73A63Ad93Ca019B25ba0aa3A28196D4",
"0x350EEA55B02150A89fe1FfFFC6044757502AD9aC",
"0xdbAE3c614a5e5FCCf21F1C769CE723E1488676d1",
"0xa90aD3abF8DEBd1e0B7B3CB98dc352E34e41aab7",
"0xf7199C6567dd1628F6192C531f7267d7b89A3732",
"0x1c64072BC1941070297790f8fD21336B2B5e4682",
"0xa3c2Ac78cb4c23800AAC0ED62f72C93d520E706f",
"0xFB6E0deeDdc4C42A8f4461c1Cad67BDd8dDA16c3",
"0x9422245375FC942F6462Dd25b76239D3290A5abA",
"0x94F04C47ae80eD6ddf41Ae0328174A5A2f3d7cb9",
"0x9a8f21A193b8bF2028d0D05fF962a787fc704644",
"0xFD61d2Bb9b38e6Ec64DE54a6C580940096F330D3",
"0x648aea5a1CDE94EEc8e84A029E50898BD01D1051",
"0x7AE41954E4285b8474e786C3112Ea52F0A0b230C",
"0xcB0BF8a574521429c3776A0484E7C52D5B5565aE",
"0xF80Efc8ed652657c16C3a920f67604824eF428B3",
"0x422F3Db577dd5853a978c62b0E9A8d1720E50e72",
"0x32c8141dCC566f454dB6Dc406d9B9548F71A8586",
"0xfec3e74ee4c0cd67be6bc851fc602d2d6885e8c3",
"0xdb7f76868689ac71238b06e497ba591c983eb10a",
"0xb8736E5ffF809146636BC2fD838682387dfB7CfF",
"0xd54E6cE783dbd74a37368D5d92263CE3Cab86FA4",
"0xcc0A9551087EC4Cc8e9E51e5b155286a4037e460",
"0x54EeCCC8F9Fac686858E7f98724ebA741853F66e",
"0xa2A0A24f04b196D5219C2E09F7ab37EaB2939CfC",
"0xa2b9107a1bCdC998C25f404914DE2043393B64Bb",
"0x36B61FF157D2EB536b45d94Ecb53802B111E3A05",
"0x3E061fAd98ed453a3437097F01f5CD6E4a8CA6fF",
"0xb4F5754D9F004685eEB709F4bC7973393940A0B2"
